# operation/rgb-tetris-wall.service
[Unit]
Description=RGB tetris wall

[Service]
Type=simple
ExecStart=/home/pi/rgb-tetris-wall/run_preproc.sh

[Install]
WantedBy=multi-user.target
